# compose/galaxy-slurm/supervisor_slurm.conf
[program:munge]
user=root
command=/usr/sbin/munged --key-file=%(ENV_MUNGE_KEY_PATH)s -F --force

[program:slurmctld]
user=root
command=/usr/sbin/slurmctld -D -L /var/log/slurm-llnl/slurmctld.log -f %(ENV_SLURM_CONF_PATH)s
autostart       = %(ENV_SLURMCTLD_AUTOSTART)s
autorestart     = true
priority        = 200

[program:slurmd]
user=root
command=/usr/sbin/slurmd -f %(ENV_SLURM_CONF_PATH)s -D -L /var/log/slurm-llnl/slurmd.log
autostart       = %(ENV_SLURMD_AUTOSTART)s
autorestart     = true
priority        = 300

